# gmond/gmond.service.in
[Unit]
Description=Ganglia Monitor Daemon
After=network-online.target

[Service]
Type=forking
PIDFile=@runstatedir@/gmond.pid
ExecStart=@sbindir@/gmond --pid-file=@runstatedir@/gmond.pid

[Install]
WantedBy=multi-user.target
