# contrib/ganglia_gmond.xml.in
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='gmond'>
  <service name='ganglia/gmond' type='service' version='1'>

     <dependency name='fs' grouping='require_all' restart_on='none' type='service'>
        <service_fmri value='svc:/system/filesystem/local'/>
     </dependency>
     <dependency name='net' grouping='require_all' restart_on='none' type='service'>
        <service_fmri value='svc:/network/physical:default'/>
     </dependency>
     <dependency name='config' grouping='require_all' restart_on='none' type='path'>
        <service_fmri value='file://localhost@sysconfdir@/gmond.conf'/>
     </dependency>

     <exec_method name='start' type='method' exec='@sbindir@/gmond' timeout_seconds='60'/>
     <exec_method name='stop' type='method' exec=':kill' timeout_seconds='60' />

     <instance name="default" enabled="true">
     </instance>

     <stability value='Unstable' />

     <template>
        <common_name>
                <loctext xml:lang='C'>Ganglia Gmon Daemon</loctext>
        </common_name>
        <documentation>
                <doc_link name='ganglia.sourceforge.net' uri='http://ganglia.sourceforge.net/' />
        </documentation>
     </template>

  </service>
</service_bundle>
