# gmetad-python/gmetad-python.service.in
[Unit]
Description=Ganglia Meta Daemon in Python
After=network-online.target

[Service]
ExecStart=@bindir@/gmetad.py -f
User=ganglia

[Install]
WantedBy=multi-user.target
