# gexec/gexec-authd/config/authd.in
#!/bin/bash
#
# chkconfig: 2345 20 80
# description: authd startup script
#
# 1. This file should be installed as /etc/init.d/authd
# 
# 2. Start authd with:
#
#      /etc/init.d/authd start
#
# Source function library.
[ -f /etc/rc.d/init.d/functions ] || exit 0
. /etc/rc.d/init.d/functions

RETVAL=0

#
# The pathname substitution in daemon command assumes prefix and
# exec_prefix are same.  This is the default, unless the user requests
# otherwise.
#
case "$1" in
    start)
        echo -n "Starting authd: "
   	    daemon @prefix@/sbin/authd -d
        RETVAL=$?
        echo
        [ $RETVAL -eq 0 ] && touch /var/lock/subsys/authd
        ;;
    stop)
        echo -n "Shutting down authd: "
        killproc authd
        RETVAL=$?
        [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/authd
        echo
        ;;
    status)
        status authd
        ;;
    restart)
        $0 stop
        $0 start
		;;
    *)
        echo "Usage: $0 {start|stop}"
        exit 1
        ;;
esac
