# gmetad/gmetad.service.in
[Unit]
Description=Ganglia Meta Daemon
After=network-online.target

[Service]
Type=forking
PIDFile=@runstatedir@/gmetad.pid
ExecStart=@sbindir@/gmetad --pid-file=@runstatedir@/gmetad.pid
EnvironmentFile=-@envdir@/gmetad

[Install]
WantedBy=multi-user.target
